# src/crontab/pg_sync
#!/bin/sh

cp /etc/shadow ~postgres/data/passwd
chown postgres:postgres ~postgres/data/passwd
