# src/crontab/purgeOMEIS
#!/bin/bash
OMEIS_BASE=%OMEIS_BASE%
OME_BIN=%OME_BIN%
PIX_PURGE=%PIX_PURGE%
find $OMEIS_BASE/Pixels -regex '.*/[0123456789]*' -atime +$PIX_PURGE -exec $OME_BIN/purge '{}' \;